# .cursor-server/data/User/History/627f08d4/hg4p.service
[Unit]
Description=Evolution API WhatsApp Service
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/evolution-api
ExecStart=/root/.nvm/versions/node/v20.10.0/bin/node /root/evolution-api/dist/main.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=evolution-api
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target 